# quadlets/gitea.container
[Unit]
Description=Gitea self-hosted Git service with LFS
After=postgres.service redis.service
Requires=postgres.service redis.service

[Container]
Image=docker.io/gitea/gitea:latest
ContainerName=gitea
Network=lab-net.network
IP=__GITEA_IP__
PodmanArgs=--add-host=postgres:__POSTGRES_IP__ --add-host=redis:__REDIS_IP__
PublishPort=3000:3000
Volume=gitea.volume:/data:Z
Environment=USER_UID=__UID__
Environment=USER_GID=__UID__
Environment=GITEA__database__DB_TYPE=postgres
Environment=GITEA__database__HOST=postgres:5432
Environment=GITEA__database__NAME=gitea
Environment=GITEA__database__USER=gitea
Environment=GITEA__server__ROOT_URL=http://__GITEA_IP__:3000/
Environment=GITEA__server__DOMAIN=__GITEA_IP__
Environment=GITEA__server__LFS_START_SERVER=true
Environment=GITEA__session__PROVIDER=redis
Environment=GITEA__session__PROVIDER_CONFIG=redis://redis:6379/0
Environment=GITEA__cache__ENABLED=true
Environment=GITEA__cache__ADAPTER=redis
Environment=GITEA__cache__HOST=redis://redis:6379/0
Environment=GITEA__security__INSTALL_LOCK=true
EnvironmentFile=%h/.local/share/podman-lab/secrets.env

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=default.target
