# quadlets/lab-java.container
[Unit]
Description=Java 21 (Eclipse Temurin) portable development environment
After=lab-net.network
Requires=lab-net.network

[Container]
Image=docker.io/library/eclipse-temurin:21-jdk
ContainerName=lab-java
Network=lab-net.network
PublishPort=8080:8080
Volume=%h/Projects:/workspace:Z
Exec=sleep infinity

[Service]
Restart=always
SuccessExitStatus=0 143

[Install]
WantedBy=default.target
