# quadlets/redis.container
[Unit]
Description=Redis 7 cache and job queue for podman-lab

[Container]
Image=docker.io/redis:7-alpine
ContainerName=redis
Network=lab-net.network
Volume=%h/.local/share/podman-lab/redis:/data:Z

[Service]
Restart=always
TimeoutStartSec=60

[Install]
WantedBy=default.target
