# quadlets/registry.container
[Unit]
Description=Podman Registry v2 — local container image storage

[Container]
Image=docker.io/registry:2
ContainerName=registry
Network=lab-net.network
IP=__REGISTRY_IP__
PublishPort=5000:5000
Volume=registry.volume:/var/lib/registry:Z
Volume=%h/.local/share/podman-lab/registry/config.yml:/etc/docker/registry/config.yml:Z,ro

[Service]
Restart=always
TimeoutStartSec=60

[Install]
WantedBy=default.target
