# quadlets/wp-agent-node.container
[Unit]
Description=Woodpecker CI agent — Node.js
After=lab-net.network woodpecker-server.service
Requires=lab-net.network

[Container]
Image=docker.io/woodpeckerci/woodpecker-agent:latest
ContainerName=wp-agent-node
Network=lab-net.network
Volume=/run/user/__UID__/podman/podman.sock:/var/run/docker.sock:Z
Environment=WOODPECKER_SERVER=woodpecker-server:9000
Environment=WOODPECKER_MAX_WORKFLOWS=1
Environment=WOODPECKER_AGENT_LABELS=type=node
EnvironmentFile=%h/.local/share/podman-lab/secrets.env

[Service]
Restart=always
TimeoutStartSec=60

[Install]
WantedBy=default.target
