# quadlets/wp-agent-rust.container
[Unit]
Description=Woodpecker CI agent — Rust
After=woodpecker-server.service

[Container]
Image=docker.io/woodpeckerci/woodpecker-agent:v3
ContainerName=wp-agent-rust
Network=lab-net.network
PodmanArgs=--add-host=woodpecker-server:10.89.1.5
Volume=/run/user/__UID__/podman/podman.sock:/var/run/docker.sock:Z
Environment=WOODPECKER_SERVER=woodpecker-server:9000
Environment=WOODPECKER_MAX_WORKFLOWS=1
Environment=WOODPECKER_AGENT_LABELS=type=rust
EnvironmentFile=%h/.local/share/podman-lab/secrets.env

[Service]
Restart=always
TimeoutStartSec=60

[Install]
WantedBy=default.target
